# backend/etc/supervisord.conf
[program:gunicorn]
command=/var/projects/skivrige/application/env/bin/gunicorn run_gunicorn:app -c etc/gunicorn.conf.py
directory=/var/projects/skivrige/application
autostart=true
atuorestart=true
stdout_logfile=/var/log/skivrige/application/stdout.log
stderr_logfile=/var/log/skivrige/application/stderr.log
